# environments/computer_rl_env/server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:xvfb]
command=Xvfb :99 -screen 0 1280x960x24 -ac +render -noreset
autorestart=true
priority=10

[program:xfce4-session]
command=xfce4-session
environment=DISPLAY=":99"
autorestart=true
priority=20
# Give desktop time to start before app
startsecs=5

[program:chrome]
command=chrome
    --no-sandbox
    --disable-dev-shm-usage
    --disable-gpu
    --disable-software-rasterizer
    --user-data-dir=/root/chrome-profile
    --remote-debugging-port=1337
    --remote-debugging-address=0.0.0.0
    --no-first-run
    --disable-default-apps
    --disable-popup-blocking
    --disable-translate
    --disable-background-networking
    --disable-sync
    --metrics-recording-only
    --disable-hang-monitor
    --disable-prompt-on-repost
    --no-default-browser-check
    --disable-client-side-phishing-detection
    --disable-component-update
    --password-store=basic
    --use-mock-keychain
    about:blank
environment=DISPLAY=":99"
autorestart=true
priority=25
# Give XFCE time to fully start before launching Chrome
startsecs=3

[program:app]
command=/app/.venv/bin/uvicorn computer_rl_env.server.app:create_app --host 0.0.0.0 --port 8000 --factory
environment=DISPLAY=":99"
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
