# environments/computer_rl_env/server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:xvfb]
command=Xvfb :99 -screen 0 1280x960x24 -ac +render -noreset
autorestart=true
priority=10

[program:xfce4-session]
command=dbus-launch --exit-with-session xfce4-session
environment=DISPLAY=":99"
autorestart=true
priority=20
# Give desktop time to start before app
startsecs=5

[program:app]
command=/app/.venv/bin/uvicorn computer_rl_env.server.app:create_app --host 0.0.0.0 --port 8000 --factory
environment=DISPLAY=":99"
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
